# supervisord/supervisor.d/04_cron_schedule.ini
[program:crond]
command=crontab /etc/cron.d/root