# supervisord/supervisor.d/04_cron_schedule.ini
[program:crond]
command=/bin/busybox crond -f
